# emacs.service
[Unit]
Description=Emacs Daemon

[Service]
User=merlin
Type=forking
ExecStart=/usr/bin/emacs --daemon
ExecStop=/usr/bin/emacsclient --eval "(kill-emacs)"
EnvironmentFile=/home/merlin/dotfiles/emacs.env
Restart=always

[Install]
WantedBy=default.target
